# etc/services.d/winbind/run
#!/usr/bin/execlineb -P
/usr/sbin/winbindd  --foreground --no-process-group "$WINBINDOPTIONS"
