# etc/services.d/smbd/run
#!/usr/bin/execlineb -P
/usr/sbin/smbd  --foreground --no-process-group "$SMBDOPTIONS"
