# etc/services.d/nmbd/run
#!/usr/bin/execlineb -P
/usr/sbin/nmbd  --foreground --no-process-group "$NMBDOPTIONS"
